# docker/MS_ADC_6000.conf
[program:v5000_6000_MS_ADC]
command = /opt/anaconda/anaconda3/bin/uwsgi --ini /home/mc5k/code/v5000_6000/MS_ADC/uwsgi.ini ;
autostart = true ;
#stopsignal=QUIT
user=root ;
startsecs = 5 ;
startretries = 1 ;
autorestart = true ;
environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v5000_6000/MS_ADC/pipeline:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmdetection" ;
redirect_stderr = true ;
stdout_logfile_maxbytes = 20MB ;
stdout_logfile = /opt/logs/backend/v5000_6000/stdout_v5000_6000.log ;
stderr_logfile = /opt/logs/backend/v5000_6000/err_v5000_6000.log ;

[program:v5000_6000_CELERY00]
command =/opt/anaconda/anaconda3/bin/celery -A config.Celery worker -l info -n local6000_0 -c 1 -P threads -Q infer_task_v5000_6000 ;
directory=/home/mc5k/code/v5000_6000/MS_ADC/ ;
autostart = true ;
user=root ;
startsecs = 5 ;
startretries = 1 ;
autorestart = true ;
environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v5000_6000/MS_ADC/pipeline:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmdetection" ;
redirect_stderr = true ;
stdout_logfile_maxbytes = 20MB ;
stdout_logfile = /opt/logs/backend/v5000_6000/c00stdout_v5000_6000.log ;
stderr_logfile = /opt/logs/backend/v5000_6000/c00err_v5000_6000 ;


[program:v5000_6000_CELERY_QUANTY]
command =/opt/anaconda/anaconda3/bin/celery -A config.Celery worker -l info -n quantifier6000 -c 1 -P threads -Q quantifier_task_v5000_6000 ;
directory=/home/mc5k/code/v5000_6000/MS_ADC/ ;
autostart = true ;
user=root ;
startsecs = 5 ;
startretries = 1 ;
autorestart = true ;
environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v5000_6000/MS_ADC/pipeline:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmdetection" ;
redirect_stderr = true ;
stdout_logfile_maxbytes = 20MB  ;
stdout_logfile = /opt/logs/backend/v5000_6000/cQstdout_v5000_6000.log ;
stderr_logfile = /opt/logs/backend/v5000_6000/cQerr_v5000_6000.log ;

[program:v5000_6000_CELERY_AUTO_INFER]
command =/opt/anaconda/anaconda3/bin/celery -A config.Celery worker -l info -n auto6000_0 -c 1 -P threads -Q auto_infer_task_v5000_6000 ;
directory=/home/mc5k/code/v5000_6000/MS_ADC/ ;
autostart = true ;
user=root ;
startsecs = 5 ;
startretries = 1 ;
autorestart = true ;
environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v5000_6000/MS_ADC/pipeline:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v5000_6000/MS_ADC/pipeline/mmlib/mmdetection" ;
redirect_stderr = true ;
stdout_logfile_maxbytes = 20MB  ;
stdout_logfile = /opt/logs/backend/v5000_6000/autostdout_v5000_6000.log ;
stderr_logfile = /opt/logs/backend/v5000_6000/autoerr_v5000_6000.log ;

#[program:CELERY_TASK_RECHECKING_STATUS_v112]
#command =/opt/anaconda/anaconda3/bin/celery -A config.Celery worker -l info -n status_0 -c 1 -P threads -Q task_rechecking_status_v112 ;
#directory=/home/mc5k/code/v1.1.2/MS_ADC/ ;
#autostart = true ;
#user=root ;
#startsecs = 5 ;
#startretries = 1 ;
#autorestart = true ;
#environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v1.1.2/MS_ADC/pipeline:/home/mc5k/code/v1.1.2/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v1.1.2/MS_ADC/pipeline/mmlib/mmdetection" ;
#redirect_stderr = true ;
#stdout_logfile_maxbytes = 20MB  ;
#stdout_logfile = /opt/logs/backend/task_rechecking_status_stdout_v112.log ;
#stderr_logfile = /opt/logs/backend/task_rechecking_status_err_v112.log ;
#
#[program:CELERY_SCHEDULE_TASK_WORKER_v112]
#command =/opt/anaconda/anaconda3/bin/celery -A config.Celery worker -l info -P eventlet ;
#directory=/home/mc5k/code/v1.1.2/MS_ADC/ ;
#autostart = true ;
#user=root ;
#startsecs = 5 ;
#startretries = 1 ;
#autorestart = true ;
#environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v1.1.2/MS_ADC/pipeline:/home/mc5k/code/v1.1.2/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v1.1.2/MS_ADC/pipeline/mmlib/mmdetection" ;
#redirect_stderr = true ;
#stdout_logfile_maxbytes = 20MB  ;
#stdout_logfile = /opt/logs/backend/schedule_worker_stdout_v112.log ;
#stderr_logfile = /opt/logs/backend/schedule_worker_err_v112.log ;
#
#[program:CELERY_SCHEDULE_TASK_v112]
#command =/opt/anaconda/anaconda3/bin/celery -A config.Celery beat -l INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler ;
#directory=/home/mc5k/code/v1.1.2/MS_ADC/ ;
#autostart = true ;
#user=root ;
#startsecs = 5 ;
#startretries = 1 ;
#autorestart = true ;
#environment=ADC_ENV=production,PYTHONPATH="$PYTHONPATH:/home/mc5k/code/v1.1.2/MS_ADC/pipeline:/home/mc5k/code/v1.1.2/MS_ADC/pipeline/mmlib/mmclassification:/home/mc5k/code/v1.1.2/MS_ADC/pipeline/mmlib/mmdetection" ;
#redirect_stderr = true ;
#stdout_logfile_maxbytes = 20MB  ;
#stdout_logfile = /opt/logs/backend/schedule_stdout_v112.log ;
#stderr_logfile = /opt/logs/backend/schedule_err_v112.log ;
